# socks-forwarder/etc/systemd/system/viam-agent.service.d/socks-forwarder.conf
[Service]
Environment="SOCKS_PROXY=localhost:1080"
